# c/Makefile/service/src/gws5k.serv
#!/bin/sh /etc/rc.common

START=90
STOP=90

APP=gws5k
CONF=/etc/config/$APP
CMD=/usr/sbin/$APP
#DEP=/usr/sbin/gws5001app

uci_conf() {
	local conf="$1"

	#if [ -f "$DEP" ]; then
		config_get key "$conf" key
	#fi
	

	$CMD -k $key >/dev/null 2>&1 &
}

start() {
	if [ -f $CONF ]; then
		config_load $CONF
		config_foreach uci_conf $APP
	else
		echo "* config file missing, using default settings"
		$CMD >/dev/null 2>&1 &
	fi
}

stop() {
	killall $APP
}
